# modules/hardware/bluetooth/files/70-bluetooth-unblock.rules
# Force unblock Bluetooth on add (Fix for stubborn cards)
ACTION=="add", KERNEL=="hci0", RUN+="/usr/bin/rfkill unblock bluetooth"
