# modules/hardware/controllers/files/99-vader-hide.rules
# Flydigi Vader 4 Pro - "Nuclear" Hiding Rule
# We physically DELETE the /dev/input nodes for the raw controller.
# Steam looks at /dev/input, so it sees nothing.
# InputPlumber looks at /dev/hidraw, so it still works.

SUBSYSTEM=="input", ATTRS{idVendor}=="04b4", ATTRS{idProduct}=="2412", ENV{ID_INPUT}:="", ENV{ID_INPUT_JOYSTICK}:="", ENV{ID_INPUT_KEY}:="", ENV{LIBINPUT_IGNORE_DEVICE}:="1", RUN+="/usr/bin/rm -f %E{DEVNAME}"
